# system/openclaw-kit-autoupdate.timer
[Unit]
Description=Periodic OpenClaw toolkit auto-update

[Timer]
OnCalendar=daily
RandomizedDelaySec=30m
Persistent=true

[Install]
WantedBy=timers.target

