# system/openclaw-kit-autoupdate.service
[Unit]
Description=Auto-update OpenClaw Codex Discord toolkit deployment
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
Group=root
Environment=OPENCLAW_KIT_REPO_DIR=/root/VibeResearch_toolkit
Environment=OPENCLAW_PROXY_ENV_FILE=/root/.openclaw/proxy.env
Environment=OPENCLAW_KIT_AUTOUPDATE_LOG=/var/log/openclaw-kit-autoupdate.log
ExecStart=/usr/local/bin/openclaw-kit-autoupdate.sh
